# systemd/arbtt.service
[Unit]
Description=arbtt service

[Service]
Type=simple
ExecStart=/home/bersam/.cabal/bin/arbtt-capture

[Install]
WantedBy=default.target
